# build-aux/deb/control/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/opentsdb" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d opentsdb stop || true
	else
		/etc/init.d/opentsdb stop || true
	fi
fi